# scripts/com.hanabi-board.weekly-backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hanabi-board.weekly-backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>/Users/yoheimizuno/hanabi-dashboard/scripts/backup.sh</string>
    </array>

    <!-- 毎週日曜 6:00 JST 実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>0</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/yoheimizuno/hanabi-dashboard/logs/backup.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yoheimizuno/hanabi-dashboard/logs/backup.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/yoheimizuno</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
